# myservice.service
[Unit]
Description=My Service

[Service]
Type=simple
ExecStart=/usr/local/bin/main.sh

[Install]
WantedBy=multi-user.target

